# cmd/supervisor.d/update_location.ini
[program:downloadlocation]
command=/var/www/flatfy/api/artisan sphinx:locations
autostart=true
autorestart=flase
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=%(here)s/../var/%(program_name)s.log
stdout_logfile=%(here)s/../var/%(program_name)s.log